# supervisord.conf
[supervisord]
nodaemon=true

[program:thin]
command=/bin/bash -c "service thin start"
